# ami/webservice.service
[Unit]
Description=index.js

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/
ExecStart=/home/ec2-user/webapp.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target